# SMBR.Services/reactor-api-server/reactor-api-server.service
[Unit]
Description=SMBR API Server Service
After=network.target
After=reactor-core-module.service
Wants=reactor-core-module.service


[Service]
Type=simple
ExecStart=/usr/bin/reactor-api-server
Restart=always
User=root

[Install]
WantedBy=multi-user.target
